# etc/systemd/system/jupyterhub.service

[Unit]
Description=Jupyterhub

[Service]
WorkingDirectory=/srv/jupyterhub
ExecStart=/usr/local/bin/jupyterhub
Restart=always

[Install]
WantedBy=multi-user.target

